# files/etc/uci-defaults/99-custom-config
#!/bin/sh

# 修改LAN接口的IP
uci set network.lan.ipaddr='192.168.0.11'

# 修改LAN接口的网关为192.168.0.1
uci set network.lan.gateway='192.168.3.158'

# 修改LAN接口的DNS为223.5.5.5
uci set network.lan.dns='223.5.5.5'


# 设置系统管理密码为root
echo -e "root\nroot" | passwd

# 主题设置
sed -i 's/bing/none/' package/openwrt-pacmluci-app-argon-config/root/etc/config/argon
curl -L https://git.kejizero.online/zhao/files/raw/branch/main/images/bg1.jpg -o package/openwrt-package/luci-theme-argon/htdocs/luci-static/argon/img/bg1.jpg

# 加载软件源
git clone --depth=1 https://github.com/oppen321/openwrt-package package/openwrt-package

echo "All done!"
